# supervisord.conf
[supervisord]
nodaemon=true

[program:shiny-server]
command=/usr/bin/shiny-server.sh --server-daemonize 0

[program:rstudio-server]
command=/usr/bin/rstudio-server.sh --server-daemonize 0

